# recipes-core/dom0-services/files/nvgpu-unbind.service
[Unit]
Description=Unbind nvgpu from Dom0 before DomD starts
Before=xendomains.service
After=sys-bus-platform-drivers-nvgpu.device

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/bin/sh -c 'echo 17000000.gpu > /sys/bus/platform/drivers/nvgpu/unbind'

[Install]
WantedBy=multi-user.target